# hosts/nuc-server/systemd/plaud-automation.service
[Unit]
Description=Run Plaud.ai to Google Drive Automation

[Service]
Type=oneshot
ExecStart=%h/github/tariqk00/plaud/bin/run_automation.sh
WorkingDirectory=%h/github/tariqk00/plaud

[Install]
WantedBy=default.target
