# services/trainheroic-extract.timer
[Unit]
Description=Run TrainHeroic Extraction Daily at 6am

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
